# docker/crontab
* * * * * cd /var/www/jobavel && php artisan schedule:run >> /dev/null 2>&1